# tower-lib/toweros-installers/toweros-thinclient/etc/local.d/00_welcome.start
if [ -f /etc/profile.d/install.sh ]; then
    echo "Welcome to TowerOS-ThinClient Installer!" > /etc/issue
    echo "" >> /etc/issue
else
    echo "Welcome to TowerOS-ThinClient!" > /etc/issue
    echo "" >> /etc/issue

    echo "Please see the ~/docs/README.md file to know how to get started with TowerOS-ThinClient." >> /etc/issue
    echo "" >> /etc/issue
fi