# tower-lib/toweros-installers/toweros-thinclient/etc/local.d/02_rfkill.start
# disable wireless devices
rfkill block all