# tower-lib/toweros-installers/toweros-thinclient/etc/local.d/04_launch_pipewire.start
# start sound server
dbus-launch /usr/libexec/pipewire-launcher >/dev/null 2>&1 &