# scripts/toweros-thinclient/etc/profile.d/install.sh
#!/bin/sh

sh /var/towercomputers/installer/install-thinclient.sh