# tower-apks/toweros-thinclient/overlay/etc/local.d/03_sync_offline_host.start
# synchronize offline hosts and Thin Client time with `router`
tower synctime