# tower-apks/toweros-thinclient/overlay/etc/local.d/03_sync_offline_host.start
# synchronize offline hosts and thin client time with `router`
tower synctime