# tower-lib/toweros-installers/toweros-thinclient/etc/init.d/supercronic
#!/sbin/openrc-run

name="supercronic"
command="supercronic"
pidfile="/run/supercronic.pid"
command_args="/etc/crontabs/supercronic"
command_background=true