# tower-apks/toweros-thinclient/overlay/etc/init.d/wl-copy-tunneler
#!/sbin/openrc-run

name="wl-copy-tunneler"
command="sh"
pidfile="/run/wl-copy-tunneler.pid"
command_args="/var/towercomputers/scripts/wl-copy-tunneler.sh"
command_background=true