# tower-apks/toweros-thinclient/overlay/etc/init.d/wl-copy-server
#!/sbin/openrc-run

name="wl-copy-server"
command="sh"
pidfile="/run/wl-copy-server.pid"
command_args="/var/towercomputers/scripts/wl-copy-server.sh"
command_background=true