# tower-apks/toweros-thinclient/overlay/etc/init.d/wl-copy-watch
#!/sbin/openrc-run

name="wl-copy-watch"
command="sh"
pidfile="/run/wl-copy-watch.pid"
command_args="/var/towercomputers/scripts/wl-copy-watch.sh"
command_background=true