# tower-apks/toweros-host/overlay/etc/init.d/xclip-copy-server
#!/sbin/openrc-run

name="xclip-copy-server"
command="sh"
pidfile="/run/xclip-copy-server.pid"
command_args="/var/towercomputers/scripts/xclip-copy-server.sh"
command_background=true
